# t410/etc/cron.d/init-tmux-session
@reboot eott /opt/microsoft/powershell/7/pwsh -File "/home/eott/scripts/try-tmux.ps1"
